# board/antc/v3-ipkvm/rootfs/etc/systemd/system/dev-binderfs.mount
[Unit]
Description=Android BinderFS
DefaultDependencies=no
[Mount]
Type=binder
What=binder
Where=/dev/binderfs
