# board/antc/v3-ipkvm/rootfs/etc/systemd/system/overlayfs.target
[Unit]
Description=OverlatFS Switch Root
DefaultDependencies=no
AllowIsolate=yes
ConditionPathIsReadWrite=!/
OnFailure=emergency.target
OnFailureJobMode=replace-irreversibly
