# board/antc/v3-ipkvm/rootfs/etc/systemd/system/overlayfs-early.target
[Unit]
Description=OverlatFS Early
DefaultDependencies=no
Before=overlayfs-ready.target
AllowIsolate=yes
ConditionPathIsReadWrite=!/
OnFailure=emergency.target
OnFailureJobMode=replace-irreversibly
