# board/antc/v3-ipkvm/rootfs/etc/systemd/system/ipmidetectd.service.d/override.conf
[Unit]
ConditionPathExists=|/etc/freeipmi/ipmidetectd.conf
ConditionPathExists=|/etc/ipmidetectd.conf
