# board/antc/v3-ipkvm/rootfs/etc/systemd/system/nfs-blkmap.service.d/override.conf
[Unit]
ConditionPathExists=/sys/module/sunrpc
ConditionPathExists=/sys/module/nfsd
