# board/antc/v3-ipkvm/rootfs/etc/systemd/system/ipmiseld.service.d/override.conf
[Unit]
ConditionPathExistsGlob=/dev/ipmi*
ConditionPathExists=|/dev/bmc
ConditionPathExists=|/dev/dcmi
