# board/antc/v3-ipkvm/rootfs/etc/systemd/system/wpa_supplicant.service.d/override.conf
[Unit]
ConditionPathExists=/sys/class/ieee80211
