# board/antc/v3-ipkvm/rootfs/etc/systemd/system/ksmbd.service.d/override.conf
[Unit]
ConditionPathExists=/etc/ksmbd/ksmbd.conf
